# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/sbin/redis-server

[program:logstash]
command=/usr/share/logstash-1.4.1/bin/logstash -f %(LOGSTASH_CONFIG_FILE) 
